# setup/stats/postfix_stats.cron
01 * * * * root /usr/bin/php -q /usr/local/issabel/postfix_stats.php
